# web2py/applications/Genios/cron/crontab
#crontab
0-59/1  *  *  *  * root *applications/Genoios/modules/exemplo.py